# worker_factory/generated_services/auto_test_worker.service
[Unit]
Description=Sameer Generated Worker auto_test_worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/sameer_ai_manager/worker_factory/generated
ExecStart=/root/sameer_ai_manager/venv/bin/python /root/sameer_ai_manager/worker_factory/generated/auto_test_worker.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
